# debian/prerm
#!/bin/sh
set -e
systemctl stop wsl-memory-guest || true
systemctl disable wsl-memory-guest || true
